# server/deploy/erythos-server.service
[Unit]
Description=Erythos sync server (Hono on Node)
Documentation=https://github.com/wolfuardian/erythos
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=erythos
Group=erythos
WorkingDirectory=/opt/erythos/server
EnvironmentFile=/opt/erythos/server/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=erythos-server

# Security hardening — see systemd.exec(5)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/erythos/server
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
